# roles/openvswitch/templates/ovs-vswitchd.service.jinja
[Unit]
After=dpdk-socket-mem.service
Requires=dpdk-socket-mem.service
